# roles/xsel/templates/xsel@.service
[Unit]
Description=Network copy to xsel

[Service]
Type=simple
ExecStart=/usr/bin/xsel -i -b -n
StandardInput=socket
StandardOutput=socket
